# spiderbridge/rootfs/etc/services.d/discovery/run
#!/bin/sh
echo "[INFO] Waiting for proxy (10s)..."
sleep 10
echo "[INFO] Starting discovery"
cd /usr/share/spiderbridge
exec python3 main_discovery.py
