# spiderbridge/rootfs/etc/services.d/mosquitto/run
#!/bin/sh
echo "[INFO] Starting mosquitto"
exec mosquitto -c /tmp/spiderbridge/mosquitto.conf
