# systemd/sf-discovery.service
[Unit]
Description=Spider Farmer Bridge — HA Discovery Publisher
After=sf-proxy.service mosquitto.service
Requires=mosquitto.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/spiderfarmer-bridge
ExecStart=/opt/spiderfarmer-bridge/.venv/bin/python3 main_discovery.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sf-discovery

[Install]
WantedBy=multi-user.target
