# systemd/sf-update.service
[Unit]
Description=Spider Farmer Bridge — Auto Update
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/spiderfarmer-bridge
ExecStart=/opt/spiderfarmer-bridge/update.sh
StandardOutput=journal
StandardError=journal
