# spiderbridge/rootfs/etc/services.d/proxy/run
#!/bin/sh
echo "[INFO] Waiting for mosquitto..."
until mosquitto_pub -h 127.0.0.1 -p 1883 -t spiderbridge/probe -m '' 2>/dev/null; do
    sleep 1
done
echo "[INFO] Starting proxy"
cd /usr/share/spiderbridge
exec python3 main_proxy.py
