# addon/spiderbridge/rootfs/etc/services.d/mosquitto/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting mosquitto"
exec mosquitto -c /tmp/spiderbridge/mosquitto.conf
